# docker/dev/worker.ini
[unix_http_server]
username=root
password=root

[supervisorctl]
username=root
password=root

[supervisord]
user=root
nodaemon=true
pidfile=/home/www-data/supervisord.pid

[program:cron]
command=/usr/sbin/crond -f
autorestart=true
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:octane]
command=php /var/www/artisan octane:start --server=swoole --host=0.0.0.0 --port=8000 --watch
autorestart=true
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:worker-default]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan queue:work --queue=default
autorestart=true
startsecs=0
stdout_logfile=/var/www/storage/logs/worker.log

[program:worker-isrc]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan queue:work --queue=isrc_lookup
autorestart=true
startsecs=0
stdout_logfile=/var/www/storage/logs/worker.log